# launchd/robinhood-codex-agent.intraday.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Replace __REPO_ROOT__ with your local checkout path before loading. -->
  <!-- Runs the package entrypoint directly. This avoids macOS launchd/TCC failures seen when
       executing repo-local shell wrappers from Documents while still keeping AGENT_ROOT explicit. -->
  <key>Label</key>
  <string>com.example.robinhood-codex-agent.intraday</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/.venv/bin/python</string>
    <string>-m</string>
    <string>trading_agent</string>
    <string>intraday</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Optional (the script also derives it), but makes the run explicit. -->
    <key>AGENT_ROOT</key>
    <string>__REPO_ROOT__</string>
    <!-- Portable PATH only. Do NOT hardcode a versioned codex release dir or a user home:
         codex auto-updates and those paths disappear, silently breaking runs. codex is
         resolved from CODEX_BIN (set a stable path in src/config/runtime.env.local if
         needed), PATH, or common install locations including ~/.local/bin. -->
    <key>PATH</key>
    <string>__REPO_ROOT__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <array>
    <!-- __INTRADAY_SCHEDULE__ -->
  </array>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/runtime/logs/launchd.intraday.out</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/runtime/logs/launchd.intraday.err</string>
</dict>
</plist>
