# launchd/robinhood-codex-agent.premarket.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Replace __REPO_ROOT__ with your local checkout path before loading. -->
  <key>Label</key>
  <string>com.example.robinhood-codex-agent.premarket</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/.venv/bin/python</string>
    <string>-m</string>
    <string>trading_agent</string>
    <string>premarket</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_ROOT</key>
    <string>__REPO_ROOT__</string>
    <key>CODEX_BIN</key>
    <string>/Users/yinyin-jiang/.codex/packages/standalone/releases/0.137.0-aarch64-apple-darwin/bin/codex</string>
    <key>PATH</key>
    <string>/Users/yinyin-jiang/.local/bin:/Users/yinyin-jiang/.codex/packages/standalone/releases/0.137.0-aarch64-apple-darwin/bin:/Users/yinyin-jiang/Documents/trade/robinhood-codex-agent/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/runtime/logs/launchd.premarket.out</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/runtime/logs/launchd.premarket.err</string>
</dict>
</plist>
