# launchd/robinhood-codex-agent.premarket.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Replace __REPO_ROOT__ with your local checkout path before loading. -->
  <key>Label</key>
  <string>com.example.robinhood-codex-agent.premarket</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/src/scripts/entrypoints/run_premarket.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/runtime/logs/launchd.premarket.out</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/runtime/logs/launchd.premarket.err</string>
</dict>
</plist>
